# agentspaces/com.claude-mem.pr-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-mem.pr-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/seb/AI/mem-claude/agentspaces/pr-watcher.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>WorkingDirectory</key>
    <string>/Users/seb/AI/mem-claude</string>

    <key>StandardOutPath</key>
    <string>/Users/seb/.claude-mem/logs/pr-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/seb/.claude-mem/logs/pr-watcher.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/seb/.local/bin:/Users/seb/.bun/bin</string>
        <key>HOME</key>
        <string>/Users/seb</string>
        <key>BEADS_NO_DAEMON</key>
        <string>1</string>
    </dict>
</dict>
</plist>
